# acarsserv.conf
description "ACARS Server"
author "Aaron Schulman"

start on runlevel [2345]

exec /usr/local/bin/acarsserv -b /acars/acarsserv.sqb -s -d -a -v >> /acars/acarsserv.log
